# etc/systemd/user/wallpaper.timer
[Unit]
Description=Change wallpaper

[Timer]
OnCalendar=*:00:00,30

[Install]
WantedBy=timers.target
