# etc/systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper

[Service]
Type=oneshot
ExecStart=/home/user/.fehbg

[Install]
WantedBy = multi-user.target
